# root/etc/udev/rules.d/80-yubikey-actions.rules
SUBSYSTEM=="usb", ACTION=="remove", ENV{ID_VENDOR_FROM_DATABASE}=="Yubico.com", RUN+="/usr/bin/systemctl --no-block start yubikey-removed.service"
