# home/.config/autostart/albert.desktop
[Desktop Entry]
Exec=/bin/albert
Name=albert
Type=Application
