# home/.config/autostart/caprine.desktop
[Desktop Entry]
Exec=/bin/caprine
Name=caprine
Type=Application
