# home/.config/autostart/chromium.desktop
[Desktop Entry]
Exec=/bin/chromium
Name=chromium
Type=Application
