# home/.config/autostart/copyq.desktop
[Desktop Entry]
Exec=/bin/copyq
Name=copyq
Type=Application
