# home/.config/autostart/discord.desktop
[Desktop Entry]
Exec=/bin/discord
Name=discord
Type=Application
