# home/.config/autostart/evolution.desktop
[Desktop Entry]
Exec=/bin/evolution
Name=evolution
Type=Application
