# home/.config/autostart/flameshot.desktop
[Desktop Entry]
Exec=/bin/flameshot
Name=flameshot
Type=Application
