# home/.config/autostart/nm-applet.desktop
[Desktop Entry]
Exec=/bin/nm-applet
Name=nm-applet
Type=Application
