# home/.config/autostart/picom.desktop
[Desktop Entry]
Exec=/bin/picom --backend glx
Name=picom
Type=Application
