# home/.config/autostart/picom.desktop
[Desktop Entry]
Exec=/bin/picom
Name=picom
Type=Application
