# home/.config/autostart/udiskie.desktop
[Desktop Entry]
Exec=/bin/udiskie --tray --no-automount
Name=udiskie
Type=Application
