# home/.config/autostart/udiskie.desktop
[Desktop Entry]
Exec=/bin/udiskie --tray
Name=udiskie
Type=Application
