# root/etc/systemd/system/yubikey-removed.service
[Unit]
Description=Yubikey Removal

[Service]
ExecStart=/usr/bin/lock-with-xsecurelock

[Install]
WantedBy=multi-user.target
